# deploy/systemd/web_speed.service
[Unit]
Description=web_speed 网络测速服务
Wants=network-online.target
After=network-online.target

[Service]
Type=exec
User=web_speed
Group=web_speed
WorkingDirectory=/var/lib/web_speed
Environment=WEB_SPEED_CONFIG=/etc/web_speed/config.toml
EnvironmentFile=-/etc/default/web_speed
ExecStart=/opt/web_speed/web_speed
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
KillSignal=SIGTERM
UMask=0077

# 服务只需要监听本地端口并写入测速数据库，默认收紧系统写入面。
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/web_speed
StateDirectory=web_speed
StateDirectoryMode=0750
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target